# config/db_failover.service
# db_failover.service
# systemd db_failover monitor unit file
# installed via daemonize script by config/db_failover

[Unit]
Description=Database failover monitoring daemon
After=network.target 


[Service]
User=ubuntu
Type=simple
ExecStart=/bin/bash -c 'cd /home/ubuntu/db_failover; python3 monitor.py'
Restart=on-failure
RestartSec=5s
PrivateTmp=True


[Install]
WantedBy=multi-user.target
